# packages/sp-raspberrypi/opt/streamplace.service
[Unit]
Description=streamplace
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/streamplace
ExecStart=/opt/streamplace/bin/streamplace
Restart=always
RestartSec=3
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
