# packages/sp-compositor/supervisor/conf.d/nwjs.conf
[program:nwjs]
command=/opt/nwjs/nw %(ENV_NWJS_ARGS)s
priority=99
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
